# wg-random.service
[Unit]
Description=Random Wireguard Connection

[Service]
Type=simple
ExecStart=/usr/local/bin/wg-random up

[Install]
WantedBy=multi-user.target
